# ucloud-geo.service
[Unit]
Description=UCloud GEO Evaluation API
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ucloud-geo-eval/backend
# 使用 xvfb-run 提供虚拟 X 显示，让 Playwright 以 headed 模式运行绕过反爬检测
ExecStart=/usr/bin/xvfb-run -a /opt/ucloud-geo-eval/venv/bin/uvicorn app:app --host 0.0.0.0 --port 8000 --workers 1
Restart=always
RestartSec=5
Environment=PYTHONPATH=/opt/ucloud-geo-eval/backend:/opt/ucloud-geo-eval/core
Environment=DISPLAY=:99

[Install]
WantedBy=multi-user.target